# conf/local/local.supervisor-api.conf
[program:API]
directory=/root/kalulu/api
command=/root/kalulu/venv/bin/gunicorn -w 4 API:app -b 127.0.0.1:5000
autostart=true
autorestart=true
stderr_logfile=/var/log/API.err.log
stdout_logfile=/var/log/API.out.log
redirect_stderr=true
